# packaging/mimir-capture/usr/lib/systemd/user/mimir-capture.service
[Unit]
Description=Mimir Activity Capture
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/mimir-capture
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
